# config/cloud/rclone-sync.timer
[Unit]
Description=Run Rclone Sync Every 3 Hours

[Timer]
OnCalendar=*:0/3
Persistent=true

[Install]
WantedBy=timers.target